# configs/systemd/5gc-nssfd.service.in
[Unit]
Description=@project_name@ NSSF Daemon
After=network-online.target

[Service]
Type=simple

User=@project_name@
Group=@project_name@

Restart=always
ExecStart=@bindir@/@project_name@-nssfd -c @sysconfdir@/@project_name@/nssf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
