# configs/systemd/5gc-mmed.service.in
[Unit]
Description=@project_name@ MME Daemon
After=network-online.target

[Service]
Type=simple

User=@5g_user@
Group=@5g_user@

Restart=always
ExecStart=@bindir@/@project_name@-mmed -c @sysconfdir@/@project_name@/mme.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID
Environment="LD_LIBRARY_PATH=@libdir@:%LD_LIBRARY_PATH%" 
[Install]
WantedBy=multi-user.target
