# configs/systemd/5gc-amf-spsd@.service.in
[Unit]
Description=5gc AMF Daemon
After=network-online.target

[Service]
Type=simple

User=@5g_user@
Group=@5g_usere@

Restart=always
ExecStart=@bindir@/@project_name@-amf-spsd -i %i -c @sysconfdir@/@project_name@/amf.yaml -l @localstatedir@/log/@project_name@/amf_sps%i.log 
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
