# configs/systemd/5gc-mmed.service.in
[Unit]
Description=@project_name@ MME Daemon
After=network-online.target

[Service]
Type=simple

User=@project_name@
Group=@project_name@

Restart=always
ExecStart=@bindir@/@project_name@-mmed -c @sysconfdir@/@project_name@/mme.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
