# configs/systemd/5gc-pcfd.service.in
[Unit]
Description=@project_name@ PCF Daemon
Wants=mongodb.service mongod.service

[Service]
Type=simple

User=@5g_user@
Group=@5g_usere@

Restart=always
ExecStart=@bindir@/@project_name@-pcfd -c @sysconfdir@/@project_name@/pcf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
