# configs/systemd/5gc-bsfd.service.in
[Unit]
Description=@project_name@ BSF Daemon
Wants=mongodb.service mongod.service

[Service]
Type=simple

User=@project_name@
Group=@project_name@

Restart=always
ExecStart=@bindir@/@project_name@-bsfd -c @sysconfdir@/@project_name@/bsf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
