# configs/systemd/5gc-pcrfd.service.in
[Unit]
Description=@project_name@ PCRF Daemon
Wants=mongodb.service mongod.service

[Service]
Type=simple

User=@5g_user@
Group=@5g_usere@

Restart=always
ExecStart=@bindir@/@project_name@-pcrfd -c @sysconfdir@/@project_name@/pcrf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
