# configs/systemd/5gc-nrfd.service.in
[Unit]
Description=@project_name@ NRF Daemon
Wants=mongodb.service mongod.service

[Service]
Type=simple

User=@5g_user@
Group=@5g_user@

Restart=always
ExecStart=@bindir@/@project_name@-nrfd -c @sysconfdir@/@project_name@/nrf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID
Environment="LD_LIBRARY_PATH=@libdir@:%LD_LIBRARY_PATH%" 
[Install]
WantedBy=multi-user.target
