# configs/systemd/5gc-udmd.service.in
[Unit]
Description=@project_name@ UDM Daemon
After=network-online.target

[Service]
Type=simple

User=@project_name@
Group=@project_name@

Restart=always
ExecStart=@bindir@/@project_name@-udmd -c @sysconfdir@/@project_name@/udm.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
