# configs/systemd/5gc-ausfd.service.in
[Unit]
Description=@project_name@ AUSF Daemon
After=network-online.target

[Service]
Type=simple

User=@project_name@
Group=@project_name@

Restart=always
ExecStart=@bindir@/@project_name@-ausfd -c @sysconfdir@/@project_name@/ausf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
