# configs/systemd/5gc-nrfd.service.in
[Unit]
Description=@project_name@ NRF Daemon
Wants=mongodb.service mongod.service

[Service]
Type=simple

User=@project_name@
Group=@project_name@

Restart=always
ExecStart=@bindir@/@project_name@-nrfd -c @sysconfdir@/@project_name@/nrf.yaml
RestartSec=2
RestartPreventExitStatus=1
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
